# docker/seldon-server-maintenance/scripts/crontab
15 0 * * * root /home/seldon/rm_old_logs.sh